# board/sbnb/sbnb/rootfs-overlay/etc/systemd/system/docker.service.d/ramdisk.conf
[Unit]
After=sbnb-storage.service
Wants=sbnb-storage.service

[Service]
Environment="DOCKER_RAMDISK=true"
