# board/reefy/reefy/rootfs-overlay/usr/lib/systemd/system/systemd-networkd-wait-online.service.d/any-interface.conf
[Service]
# 1s timeout: we don't need network for boot. Docker uses cached images
# (--pull missing), reconciler applies saved state offline. Network is
# picked up organically when available (WiFi, ethernet, or never).
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any --timeout=1
