# board/reefy/reefy/rootfs-overlay/etc/systemd/system/docker.service.d/cpu.conf
[Service]
# Reserve 10% of CPU for system services (watchdog, MQTT, SSH).
# Without this, a container fork-bomb pegs all cores → systemd can't
# pet the hardware watchdog → device reboots after 120s.
ExecStartPost=/usr/bin/reefy-docker-cpuquota
