# board/reefy/reefy/rootfs-overlay/etc/systemd/system/docker.service.d/memory.conf
[Service]
# Cap Docker + all containers at 90% of system RAM (matches CPUQuota).
# 10% reserve ensures system services (reefy-mqtt, SSH, watchdog,
# systemd, kernel) have memory headroom even during container OOM
# kill/respawn storms.
MemoryMax=90%
