# board/reefy/reefy/rootfs-overlay/etc/systemd/system/docker.service.d/ramdisk.conf
[Unit]
After=reefy-storage.service
Wants=reefy-storage.service

[Service]
Environment="DOCKER_RAMDISK=true"
