# board/reefy/reefy/rootfs-overlay/etc/systemd/system/docker.service.d/stop-timeout.conf
[Service]
# systemd's default 90s eats most of the firmware-update reboot
# budget when one container drains slowly (cloudflared in particular
# waits for graceful tunnel shutdown). 30s is enough for our typical
# workload - containers stop in <10s, backups use LVM snapshots so
# crash-consistent shutdown is fine. Same speedup applies to manual
# reboots, not just OTA.
TimeoutStopSec=30s
