# board/sbnb/sbnb/rootfs-overlay/usr/lib/systemd/system/sbnb-boot-watchdog.service
[Unit]
Description=sbnb boot watchdog (A/B update safety)
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/sbnb-boot-watchdog

[Install]
WantedBy=multi-user.target
