# board/sbnb/sbnb/rootfs-overlay/usr/lib/systemd/system/sbnb-init.service
[Unit]
Description=sbnb device init (credentials, banner)
After=sbnb-storage.service
Wants=sbnb-storage.service

[Service]
Type=oneshot
ExecStart=/usr/bin/boot-sbnb-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
