# board/reefy/reefy/rootfs-overlay/usr/lib/systemd/system/reefy-mqtt.service
[Unit]
Description=Reefy MQTT Reconciler
Documentation=https://github.com/reefyai/reefy

# Start after storage setup + hostname. Hostname feeds the MQTT
# client-id and topic prefix — without it we'd register as "buildroot"
# until reefy-hostname.service eventually ran, flapping the device's
# identity on the backend. No network-online.target dependency —
# reconciler applies saved state offline and connects to MQTT when
# network becomes available.
After=reefy-storage.service reefy-hostname.service
Wants=reefy-storage.service reefy-hostname.service

# No ConditionPathExists — reconciler starts unconditionally
# and waits internally for config/certs to appear

[Service]
Type=simple

# Load configuration if exists (USB first, persistent override second)
EnvironmentFile=-/mnt/reefy/mqtt/mqtt.conf
EnvironmentFile=-/mnt/reefy-data/state/mqtt.conf

# PYTHONUNBUFFERED ensures print output appears immediately in journal
Environment=PYTHONUNBUFFERED=1
Environment=LANG=en_US.UTF-8
Environment=LC_ALL=en_US.UTF-8
ExecStart=/usr/bin/reefy-mqtt-reconciler

# Always restart — MQTT is the primary communication channel
Restart=always
RestartSec=10

# Resource limits. Soft knobs only: CPUWeight yields to user apps
# under contention, MemoryHigh throttles page reclaim above 1G. Hard
# caps (CPUQuota / MemoryMax) would kill or stall spawned children
# like `borg list`/`borg extract` mid-restore on large archives.
MemoryHigh=1G
CPUWeight=20

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reefy-mqtt

[Install]
WantedBy=multi-user.target
