# board/reefy/reefy/rootfs-overlay/usr/lib/systemd/system/reefy-boot-watchdog.service
[Unit]
Description=reefy boot watchdog (A/B update safety)
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/reefy-boot-watchdog

[Install]
WantedBy=multi-user.target
