# board/reefy/reefy/rootfs-overlay/usr/lib/systemd/system/reefy-ttyd.service
[Unit]
Description=reefy web terminal (ttyd)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ttyd -W -i 127.0.0.1 -p 7681 tmux new-session -A -s host /bin/bash -l
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
