# board/reefy/reefy/rootfs-overlay/usr/lib/systemd/system/reefy-watchdog.service
[Unit]
Description=reefy watchdog

[Service]
Type=oneshot
ExecStart=/usr/bin/reefy-watchdog.sh
TimeoutStartSec=50s
